# helloDaemon/files/helloDaemon.init
#!/bin/sh /etc/rc.common

START=98

start() {
	service_start /bin/helloDaemon
}

stop() {
	service_stop /sbin/helloDaemon
}
